# vm300_files/etc/init.d/serial_enable
#!/bin/sh /etc/rc.common

START=10
STOP=15

start()
{
	echo -n > /dev/ttyS0
}
